# install/manifests/init-lockserver.service
[Unit]
#After=

[Service]
ExecStart=/usr/local/bin/init-lockserver.sh

[Install]
WantedBy=default.target