# install/manifests/etc/systemd/init-protheus.service
[Unit]
#After=

[Service]
ExecStart=/usr/local/bin/init-protheus.sh

[Install]
WantedBy=default.target