# install/manifests/etc/systemd/init-protheus-rest.service
[Unit]
#After=

[Service]
ExecStart=/usr/local/bin/init-protheus-rest.sh

[Install]
WantedBy=default.target