# rootfs/machine/etc/systemd/init-dbaccess.service
[Unit]
#After=

[Service]
ExecStart=/usr/local/bin/init-dbaccess.sh

[Install]
WantedBy=default.target